# scripts/com.jdot.apple-mail-autosort.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd job: run the local Gemma-driven inbox sorter once a day.

  BEFORE LOADING, edit the two paths marked TODO:
    • the node binary  (run `which node` in Terminal — on Apple-Silicon Homebrew
      it's usually /opt/homebrew/bin/node)
    • the absolute path to auto-sort.mjs

  Install:
    cp scripts/com.jdot.apple-mail-autosort.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.jdot.apple-mail-autosort.plist
  Disable:
    launchctl unload ~/Library/LaunchAgents/com.jdot.apple-mail-autosort.plist

  Note: this is set to DRY_RUN=0 (it will actually move mail). Test by hand first
  (see AUTOSORT.md) so macOS grants the Automation permission to control Mail.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jdot.apple-mail-autosort</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>                                   <!-- TODO: `which node` -->
    <string>/Users/jdot/Documents/Development/apple-mail-mcp/scripts/auto-sort.mjs</string>  <!-- TODO: confirm path -->
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- On-device by default. To use a model on another machine you must both
         point this at it AND add  <key>AUTOSORT_ALLOW_REMOTE</key><string>1</string>
         below — otherwise the script refuses to send mail metadata off-device. -->
    <key>LMSTUDIO_URL</key>   <string>http://localhost:1234</string>
    <!-- Pin a model id if you don't want auto-detect:
    <key>LMSTUDIO_MODEL</key> <string>gemma-3-4b-it</string> -->
    <key>ACCOUNTS</key>      <string>iCloud,Google</string>
    <!-- Ships as a dry run. Flip to "0" ONLY after a manual dry run looks right. -->
    <key>DRY_RUN</key>       <string>1</string>
  </dict>

  <!-- Run daily at 07:00 local time. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>7</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <!-- The script writes its own owner-only (0600) log to
       ~/Library/Logs/apple-mail-autosort.log. We deliberately do NOT redirect
       stdout/stderr into world-readable /tmp. If you want launchd's raw
       stdout too, point these at an absolute path under your own ~/Library/Logs
       (TODO: replace <YOU> with your username) — never /tmp on a shared Mac. -->
  <key>RunAtLoad</key>       <false/>
  <key>StandardOutPath</key> <string>/Users/&lt;YOU&gt;/Library/Logs/apple-mail-autosort.launchd.log</string>
  <key>StandardErrorPath</key><string>/Users/&lt;YOU&gt;/Library/Logs/apple-mail-autosort.launchd.log</string>
</dict>
</plist>
